# packer/files/puma.service
[Unit]
Description=puma
After=network.target

[Service]
Type=forking
User=appuser
Group=appuser
ExecStart=/usr/local/bin/puma --dir /home/appuser/reddit/ --pidfile /tmp/puma.pid --daemon
ExecStop=/bin/kill -HUP $MAINPID
PIDFile=/tmp/puma.pid
Restart=no

[Install]
WantedBy=multi-user.target

